# root/etc/systemd/system/evebox.service.d/eve-json.conf
[Unit]
ConditionPathExists=/var/log/suricata/eve.json
After=suricata.service

[Service]
User=suricata
Group=suricata

